# templates/s3.upstart.conf
description "Mount Amazon S3 file system on system start" 

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec s3fs -f <%= config["s3_bucket"] %> /mnt/s3 -o use_cache=/tmp -o allow_other